# config-generator.ini
[inet_http_server]
port = 127.0.0.1:9001

[program:config-generator]
command=python3 /usr/local/bin/generate_config
startretries=0
autostart=true
startsecs = 0
autorestart=false
redirect_stderr=true
exitcodes=1
